# meta-opencentauri/recipes-core/zram/files/zram-emmc-swap
#!/bin/sh
### BEGIN INIT INFO
# Provides:          zram-emmc-swap
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: eMMC swapfile fallback for low-memory systems
### END INIT INFO

SWAP_FILE="/mnt/writable/swapfile"

get_swap_size_mb() {
    totalmem_kb=$(awk '/^MemTotal/ {print $2}' /proc/meminfo)
    echo $((totalmem_kb / 1024 * 2))
}

case "$1" in
    start)
        if [ ! -d /mnt/writable ]; then
            echo "zram-emmc-swap: /mnt/writable not mounted, skipping"
            exit 0
        fi

        swap_mb=$(get_swap_size_mb)

        if [ ! -f "$SWAP_FILE" ]; then
            echo "zram-emmc-swap: creating ${swap_mb}MB swapfile..."
            dd if=/dev/zero of="$SWAP_FILE" bs=1M count="$swap_mb" 2>/dev/null
            chmod 0600 "$SWAP_FILE"
            mkswap "$SWAP_FILE"
        elif [ "$(stat -c%s "$SWAP_FILE" 2>/dev/null)" != "$((swap_mb * 1024 * 1024))" ]; then
            echo "zram-emmc-swap: swapfile size mismatch, recreating..."
            swapoff "$SWAP_FILE" 2>/dev/null
            rm -f "$SWAP_FILE"
            dd if=/dev/zero of="$SWAP_FILE" bs=1M count="$swap_mb" 2>/dev/null
            chmod 0600 "$SWAP_FILE"
            mkswap "$SWAP_FILE"
        fi

        swapon -p 10 "$SWAP_FILE" 2>/dev/null
        ;;
    stop)
        swapoff "$SWAP_FILE" 2>/dev/null
        ;;
    status)
        if grep -q "$SWAP_FILE" /proc/swaps 2>/dev/null; then
            echo "zram-emmc-swap: active"
            grep "$SWAP_FILE" /proc/swaps
        else
            echo "zram-emmc-swap: inactive"
        fi
        ;;
    *)
        echo "Usage: $0 {start|stop|status}"
        exit 1
        ;;
esac

exit 0
