# meta-opencentauri/recipes-core/busybox/files/crontab
*/5 * * * * /usr/sbin/logrotate /etc/logrotate.conf
*/15 * * * * /usr/sbin/ntpd -nq -p pool.ntp.org
